# deploy/adventure-rpg.service
# /etc/systemd/system/adventure-rpg.service
# Next.js 프로덕션 서버를 항상 가동 + 죽으면 자동 재시작 + 서버 부팅 시 자동 시작.
# 환경변수는 .env.production.local 에 두면 `next start` 가 알아서 로드한다 (EnvironmentFile 불필요).
[Unit]
Description=Adventure RPG
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ec2-user/adventure-rpg
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3
User=ec2-user
Group=ec2-user
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
